# invoiceplane/rootfs/etc/service/ip/run
#!/bin/sh
set -e

document_root="{{ document_root }}"

# Preveri, ali je mapa $document_root že nameščena
if [ -d "$document_root" ]; then
  echo "InvoicePlane je že nameščen."
  ln -s $document_root/ip/* /var/www/html/
  echo "Mapa uspešno povezana na strežnik"
  exit 1
fi

# Preveri, ali je dosegljiva povezava do GitHuba
if ! curl --output /dev/null --silent --head --fail "https://github.com/InvoicePlane/InvoicePlane/releases/download/v1.6.1-beta-2/v1.6.1-beta-2.zip"; then
  echo "Ni mogoče prenesti InvoicePlane. Povezava do GitHuba ni na voljo."
  exit 1
fi

# Ustvari mapo, če ta ne obstaja
mkdir -p $document_root

# Prenesi InvoicePlane
curl -LJO https://github.com/InvoicePlane/InvoicePlane/releases/download/v1.6.1-beta-2/v1.6.1-beta-2.zip

# Razpakiraj arhiv v mapo /share/www/invoiceplane
unzip v1.6.1-beta-2.zip -d $document_root

# Odstrani preneseni arhiv
rm -f v1.6.1-beta-2.zip

# Sporoči, da je bila namestitev InvoicePlane uspešna
echo "InvoicePlane je bil uspešno nameščen v $document_root."
