# nginx_n_php/rootfs/etc/s6-overlay/s6-rc.d/www/run

#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Start php service
# ==============================================================================

set -e

# connect shaer folder with web server
ln -s /share/html /var/www/
