# rpi-camera-rtsp-server-onvif/data/addon/rpi-camera-rtsp-server-onvif/etc/services.d/rpos/run
#!/usr/bin/with-contenv bashio
echo "$( date +'%Y/%m/%d %H:%M:%S' ) Starting Rpos"
cd /data/rpos
exec  /usr/local/bin/node /data/rpos/rpos.js
